# ansible/templates/update-motd.d/98-hashistack.j2
#!/bin/sh
printf "Hashistack versions:"
printf "Consul:\t\t$(consul -version)"
printf "Vault:\t\t$(vault -version)"
printf "Terraform:\t$(terraform -version)"
printf "Nomad:\t\t$(nomad -version)"