# ansible/templates/update-motd.d/98-hashistack.j2
#!/bin/sh
echo "Hashistack versions:"
echo "Consul:\t\t$(consul -version)"
echo "Vault:\t\t$(vault -version)"
echo "Terraform:\t$(terraform -version)"
echo "Nomad:\t\t$(nomad -version)"