# ansible/templates/update-motd.d/99-help-text.j2
#!/bin/sh
echo "\n## NOTE ##"
echo "If you started from scratch without a template see: https://github.com/fredrikhgrelland/vagrant-hashistack/blob/master/README.md#usage"
echo "##########"
