# ansible/templates/update-motd.d/99-help-text.j2
#!/bin/sh
printf "\n## NOTE ##"
printf "If you started from scratch without a template see: https://github.com/fredrikhgrelland/vagrant-hashistack/blob/master/README.md#usage"
printf "##########"
